# scripts/systemd/signal-bot-dashboard.service
[Unit]
Description=Signal Bot Dashboard
After=network.target signal-bot.service

[Service]
Type=simple
User=zknowles
WorkingDirectory=/home/zknowles/signal-bot/dashboard
ExecStart=/usr/bin/npx tsx src/server.ts
Restart=on-failure
RestartSec=5
Environment=DASHBOARD_PORT=3333
Environment=FACTORY_RUNS_DIR=/home/zknowles/signal-bot/factory/runs

[Install]
WantedBy=multi-user.target
